# 启动项插件/ping
#!/bin/sh /etc/rc.common

START=49
start()
{
	echo start
	/test/ping.sh &
}

